# meta-vortex/meta-vortex-var-som-mx6/recipes-vortex/vortex-init/files/vortex-init.service
[Unit]
Description="Initialisation script for CAN and IOs."

[Service]
Type=oneshot
ExecStart=/usr/bin/vortex-init.sh

[Install]
WantedBy=multi-user.target

