# meta-vortex/meta-vortex-var-som-mx6/recipes-vortex/vortex-usb-wan/files/iphone-pairing.rules
SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ENV{PRODUCT}=="5ac/12[9a][0-9a-f]/*|5ac/8600/*", ACTION=="add",  , RUN+="/usr/bin/iphone-pairing.sh"
SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ENV{PRODUCT}=="5ac/12[9a][0-9a-f]/*|5ac/8600/*", ACTION=="remove", RUN+="/usr/bin/iphone-pairing.sh"
