# meta-vortex/meta-vortex-var-som-mx6/recipes-connectivity/ppp/files/pppd.path
[Unit]
Description=Monitor USB serial on cell modem

[Path]
PathChanged=/dev/ttyUSB3
Unit=pppd.service

[Install]
WantedBy=multi-user.target