# meta-vortex/meta-vortex-var-som-mx6/recipes-vortex/quectel/files/pppd.service
[Unit]
Description=PPP over Serial link
After=network.target

[Service]
ExecStart=/usr/sbin/pppd -detach call quectel-ppp
Restart=always

[Install]
WantedBy=multi-user.target